# config/listen.conf
#!upstart
description "Listener"

start on runlevel [2345]
stop on runlevel [!2345]

exec /home/ubuntu/cfm/manage.py listen asdf

respawn

# Give up if restart occurs 10 times in 90 seconds.
respawn limit 10 90