# config/listen.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /home/ubuntu/cfm/manage.py listen test